# pixie-media/live-build/config/includes.chroot/etc/systemd/system/getty@.service.d/zzz-pixie-autologin-root.conf
# pixie live env: auto-login root on the alternate VTs (tty2-tty6).
#
# live-config's stock setup tries to ``agetty --autologin user`` on
# every VT, where ``user`` is the account live-config-systemd creates
# at first boot. In our chroot that user-creation path is racing /
# silently failing (the user has reported an "Authentication failure"
# loop on tty2 with the autologin retrying forever) and the operator
# loses the Alt+F2 diagnostic shell.
#
# We override at the template level so it applies to every getty@
# instance the system spawns on demand. tty1 is unaffected at runtime
# because ``pixie-on-tty1.service`` declares ``Conflicts=getty@tty1.
# service`` -- the getty there never starts in the first place, so
# this drop-in's autologin-root never runs on tty1.
#
# Root autologin is acceptable here because the pixie live env is a
# single-purpose flashing appliance: anyone with physical Alt+Fn
# access already has full local access via the boot menu (they can
# add ``init=/bin/sh`` and bypass everything). The convenience of an
# instant root shell for diagnostics outweighs the negligible delta.
#
# ``zzz-`` filename prefix sorts after live-config's drop-in
# (typically named ``live-config.conf`` or similar) so this entry
# wins when systemd merges the directory.

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
